# deploy/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/dev/stdout
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

# Flask AI 偵測服務（5000）。啟動時會載入模型，較慢屬正常。
[program:flask]
command=python app1.py
directory=/app/backend/python
autorestart=true
startsecs=0
environment=PYTHONIOENCODING="utf-8",HF_HUB_OFFLINE="1",TRANSFORMERS_OFFLINE="1"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# Next.js（3000）
[program:next]
command=npm run start
directory=/app/backend
autorestart=true
startsecs=0
environment=PORT="3000",HOSTNAME="0.0.0.0"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# nginx（7860 對外）
[program:nginx]
command=nginx -g 'daemon off;'
autorestart=true
startsecs=0
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
